# grafana/rootfs/etc/services.d/memcached/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Community Hass.io Add-ons: Grafana
# Runs the memcached Server
# ==============================================================================
declare -a options

bashio::log.info 'Starting Memcached...'

options+=(-l 127.0.0.1)
options+=(-m 32)

# Run Memcached
exec s6-setuidgid memcache /usr/bin/memcached "${options[@]}"
