# grafana/rootfs/etc/services.d/grafana/run
#!/usr/bin/with-contenv bashio
# shellcheck disable=SC2191
# ==============================================================================
# Community Hass.io Add-ons: Grafana
# Runs the Grafana Server
# ==============================================================================
declare -a options
declare name
declare value

# Wait for Memcached to become available
bashio::net.wait_for 11211

bashio::log.info 'Starting Grafana...'

options+=(--config="/etc/grafana/grafana.ini")
options+=(--homepath="/usr/share/grafana")

# Load custom environment variables
for var in $(bashio::config 'env_vars|keys'); do
    name=$(bashio::config "env_vars[${var}].name")
    value=$(bashio::config "env_vars[${var}].value")
    bashio::log.info "Setting ${name} to ${value}"
    export "${name}=${value}"
done

# Run Grafana
exec grafana-server "${options[@]}"
