# grafana/rootfs/etc/s6-overlay/s6-rc.d/nginx/run
#!/command/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: grafana
# Runs the Nginx daemon
# ==============================================================================

# Wait for grafana to become available
bashio::net.wait_for 3000

bashio::log.info "Starting NGINX..."

exec nginx
